# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/usr/sbin/mysqld --bind-address=0.0.0.0
autostart=true
autorestart=true

[program:http_server]
command=python3 /http_server.py
autostart=true
autorestart=true